# src/etc/local.d/change_ownership_of_dot_distcc_directory.start
#!/bin/bash

sudo chmod -R -s /var/tmp/portage/.distcc
sudo chown -R portage:pi  /var/tmp/portage/.distcc
