# src/etc/local.d/mount_gluster_volume.start
#!/bin/bash

if [[ -n /var/gluster/logs ]] ; then
   mkdir -p /var/gluster/logs ;

fi
/usr/bin/sudo /bin/mount -t glusterfs -o log-file=/var/log/gluster.log pi2:logs /var/gluster/logs

if [[ -n /var/gluster/capture ]] ; then
   mkdir -p /var/gluster/capture ;
fi
/usr/bin/sudo /bin/mount -t glusterfs -o log-file=/var/log/gluster.log pi2:capture /var/gluster/capture
